# rootfs/etc/supervisor/conf.d/000_cron.conf
[program:cron]
command=cron -fL15 
autostart=true
autorestart=true
stdout_logfile=/var/log/docker-out.log
stderr_logfile=/var/log/docker-out.log
