# rootfs/etc/supervisor/conf.d/000_cron.conf
[program:cron]
user=root
command=cron -fL15 
autostart=true
autorestart=unexpected
exitcodes=0,2
stopsignal=TERM
stdout_logfile=NONE
redirect_stderr=true
